# user-session.service.in

#
# Minimal user session launcher for systemd.
#

[Unit]
Description=User Session launcher
Wants=syslog.target dbus.service 

[Service]
ExecStart=@prefix@/bin/user-session
Restart=always
RestartSec=10

[Install]
Alias=display-manager.service
WantedBy=graphical.target
