# xorg.service.in

#
# Minimal Xorg service file - launches Xorg as a service unit
#

# The Xorg launch helper forks, launches Xorg and waits for Xorg to
# accept incoming connections to $DISPLAY, and then signals READY
# to systemd. This guarantees that services that require access to
# $DISPLAY during the session don't start too early.
#
# If you implement a service that requires access to $DISPLAY, your
# service unit file needs to include 'After=xorg.target'.

[Unit]
Description=Xorg server launch helper
Before=xorg.target
After=getty@tty1.service rc-local.service plymouth-start.service systemd-user-sessions.service plymouth-quit.service livesys-late.service
Conflicts=getty@tty1.service plymouth-quit.service
OnFailure=plymouth-quit.service

[Service]
Type=notify
ExecStart=@bindir@/xorg-launch-helper
Restart=always
RestartSec=10

[Install]
WantedBy=xorg.target
